# infrastructure/nextcloud/backup/timerbackup.service
[Unit]
Description=Borg backup do Nextcloud

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/backupnextcloud
